# tldr-bot.service
[Unit]
Description=Job that runs the tldr-bot

[Service]
EnvironmentFile=${APP_ROOT}/.bootstraprc
ExecStart=/usr/local/bin/uwsgi --socket localhost:${PORT} --wsgi-file ${APP_ROOT}/post_comment.py --callable app --single-interpreter
RuntimeDirectory=uwsgi
Restart=always
KillSignal=SIGQUIT
Type=notify
StandardError=syslog
NotifyAccess=all

[Install]
WantedBy=multi-user.target
