# tldr-bot.service
[Unit]
Description=Job that runs the tldr-bot

[Service]
EnvironmentFile=/home/agniva/.bootstraprc
ExecStart=/usr/local/bin/uwsgi --socket localhost:${PORT} --wsgi-file /home/agniva/post_comment.py --callable app --master
RuntimeDirectory=uwsgi
Restart=always
KillSignal=SIGQUIT
Type=notify
StandardError=syslog
NotifyAccess=all

[Install]
WantedBy=multi-user.target
